# tradingagents-cn/rootfs/etc/services.d/tradingagents/finish
#!/bin/bash

echo "TradingAgents-CN service stopped"
